# user-config/dot-config/systemd/user/battery-notify.timer
[Unit]
Description=Battery level warning timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
