# user-config/dot-config/systemd/user/battery-notify.service
[Unit]
Description=Notify when battery level changes beyond thresholds

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-notify.sh

[Install]
WantedBy=default.target
